# core/debian/mailman-i18n.prerm
#!/bin/sh -e
# $URL$
# $Id$

if [ "$1" = "upgrade" ] ; then
  if [ -x "/etc/init.d/mailman" ]; then
    if [ -x /usr/sbin/invoke-rc.d ] ; then
      invoke-rc.d mailman stop || true
    else
      /etc/init.d/mailman stop || true
    fi
  fi
  chmod -x /var/lib/mailman/bin/mailmanctl
fi
